# update-motd.d/031-system-load
#!/bin/bash
#
# Copyright (C) TEONITE - http://teonite.com
# Copyright (C) Robert Olejnik <robert@teonite.com>
#


U=`uptime | awk '{print $10}' | sed -e 's/,//g'`
printf "  System load: %6s" ${U}
