# debian/update-motd.cron.d
5 * * * * root /usr/bin/update-motd &> /dev/null
